# conf/poller.ini
[program:poller]
command=/opt/jianyan/application/poller -config /opt/jianyan/application/config.yml -remote http://3.10sh.cn:4321/api/v1/app_settings/dump

directory=/opt/jianyan/application
process_name=%(program_name)s
user=chenwei
priority=999
autostart=true
startsecs=1
startretries=8
autorestart=true
stdout_logfile=/var/log/%(program_name)s.out.log
stderr_logfile=/var/log/%(program_name)s.err.log
